# fluidsynth.service.in
[Unit]
Description=FluidSynth Daemon
Documentation=man:fluidsynth(1)
After=sound.target
After=pipewire.service pulseaudio.service
Wants=pipewire.service pulseaudio.service
# If you need more than one instance, use `systemctl edit` to override this:
ConditionPathExists=!/run/lock/fluidsynth.lock
ConditionUser=!@system

[Service]
# added automatically, for details please see
# https://en.opensuse.org/openSUSE:Security_Features#Systemd_hardening_effort
ProtectSystem=full
ProtectHome=read-only
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
# end of automatic additions
# required in order for the above sandboxing options to work on a user unit
PrivateUsers=yes
Type=notify
NotifyAccess=main
Environment=OTHER_OPTS= SOUND_FONT=
EnvironmentFile=@FLUID_DAEMON_ENV_FILE@
EnvironmentFile=-%h/.config/fluidsynth
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/fluidsynth -is $OTHER_OPTS $SOUND_FONT
ExecStartPre=touch /run/lock/fluidsynth.lock
ExecStopPost=rm -f /run/lock/fluidsynth.lock

[Install]
WantedBy=default.target

